# usage-store.service
[Unit]
Description=usage-store — token + subscription-limit + per-key spend aggregator
After=network.target

[Service]
ExecStart=/home/kayushkincom/bin/usage-store-server
Restart=on-failure
RestartSec=3
Environment=USAGE_STORE_LISTEN_ADDR=:8185
# Spend collection is disabled until AUTH_STORE_TOKEN is set in a host-local
# drop-in (admin keys are then auto-discovered from auth-store by api_key
# prefix — no per-provider config needed). Recommended path:
#   ~/.config/systemd/user/usage-store.service.d/spend.conf

[Install]
WantedBy=default.target
